# ACS_v.18.20.4_1/ACS/acs_test_scripts/Device/Model/LinuxDevice/embedded/command_server.service
[Unit]
Description=bluetooth control service for test automation
After=network.target

[Service]
Type=simple
#EnvironmentFile=/etc/conf.d/squeezelite
ExecStart=/usr/bin/python /usr/local/bin/command_server.py
RestartSec=5
Restart=on-failure

[Install]
WantedBy=
